# supervisor-php/conf/10-url-crawler.conf
[program:url-crawler]
process_name=%(program_name)s_%(process_num)02d
directory=/srv/application/
command=php url-crawl.php
autostart=true
autorestart=true
user=root
numprocs=1